# confs/gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/django_project
command=/home/ubuntu/django_project/env/bin/gunicorn --workers 4 --bind unix:/home/ubuntu/django_project/app.sock speech_analytics_platform.wsgi:application -t 7200 --log-level debug
autostart=true
autorestart=true

stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[group:guni]
programs:gunicorn

